# .config/systemd/user/nm-applet.service
[Unit]
Documentation=man:nm-applet(1)
Description=NetworkManager Applet
PartOf=graphical-session.target

[Service]
ExecStart=nm-applet
