# .config/systemd/user/picom.service
[Unit]
Description=Picom compositor
Documentation=man:picom(1)
PartOf=graphical-session.target

[Service]
ExecStart=picom
