# .config/systemd/user/polybar.service
[Unit]
Description=Polybar status bar
Documentation=man:polybar(1)
PartOf=graphical-session.target

[Service]
ExecStart=polybar -r
